# plugins/grpc/build/config/supervisord/grpc_server.conf
[program:grpc_server]
command = kaldi_serve_app /home/app
startsecs = 0
autorestart = false
startretries = 1
redirect_stderr = true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
